# meta-application/meta-demo/recipes-helloworld/service/service-0.1/test_systemd.service
[Unit]
Description=helloworld test
Documentation=none

[Service]
ExecStart=/home/root/test_systemd.sh
Type=simple
KillMode=process
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
