# meta-core/meta-security/meta-secure-storage/recipes-utils/libsecstore/files/secstore-manager.service
[Unit]
Description=secstore-manager

[Service]
ExecStart=/usr/bin/secstore-manager start
Type=simple
KillMode=process

[Install]
WantedBy=sysinit.target
